# deploy/gauranga.service
[Unit]
Description=GAURANGA AI Agent Server
After=network.target

[Service]
Type=simple
User=openhands
WorkingDirectory=/workspace/project/Alpha-agent-Gaurangga/server
ExecStart=/workspace/project/Alpha-agent-Gaurangga/server/venv/bin/python gauranga_server.py
Restart=always
RestartSec=10
StandardOutput=append:/workspace/project/Alpha-agent-Gaurangga/logs/gauranga.log
StandardError=append:/workspace/project/Alpha-agent-Gaurangga/logs/gauranga_error.log

# Environment variables
EnvironmentFile=/workspace/project/Alpha-agent-Gaurangga/server/.env

[Install]
WantedBy=multi-user.target
